# scripts/systemd/ironclaw-sample-agent.service
# Ironclaw sample-agent — start at boot (systemd)
# Install: sudo cp scripts/systemd/ironclaw-sample-agent.service /etc/systemd/system/
#         sudo systemctl daemon-reload && sudo systemctl enable ironclaw-sample-agent.service
# Edit User, Group, WorkingDirectory, and ExecStart* to your user and repo path (e.g. on Pi: User=pi, Group=pi, /home/pi/ironclaw).

[Unit]
Description=Ironclaw sample-agent at boot
After=docker.service network-online.target
Wants=docker.service network-online.target

[Service]
Type=oneshot
User=kosar
Group=kosar
WorkingDirectory=/path/to/ironclaw
ExecStartPre=/usr/bin/test -f /path/to/ironclaw/agents/sample-agent/agent.conf
ExecStartPre=/usr/bin/test -f /path/to/ironclaw/agents/sample-agent/.env
ExecStart=/path/to/ironclaw/scripts/start-sample-agent-at-boot.sh
TimeoutStartSec=600
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
