# scripts/systemd/ironclaw-pibot.service
# Ironclaw pibot — start at boot (systemd)
# Install: sudo cp scripts/systemd/ironclaw-pibot.service /etc/systemd/system/
#         sudo systemctl daemon-reload && sudo systemctl enable ironclaw-pibot.service
# Edit User, Group, WorkingDirectory, and ExecStart* to your user and repo path (e.g. on Pi: User=pi, Group=pi, /home/pi/ironclaw). Agent name pibot is the Pi-style example; replace with your agent name if different.

[Unit]
Description=Ironclaw pibot agent at boot
After=docker.service network-online.target
Wants=docker.service network-online.target

[Service]
Type=oneshot
User=kosar
Group=kosar
WorkingDirectory=/path/to/ironclaw
# Pre-flight: fail if agent or .env missing (replace path if repo is elsewhere)
ExecStartPre=/usr/bin/test -f /path/to/ironclaw/agents/pibot/agent.conf
ExecStartPre=/usr/bin/test -f /path/to/ironclaw/agents/pibot/.env
ExecStart=/path/to/ironclaw/scripts/start-pibot-at-boot.sh
TimeoutStartSec=600
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
