# scripts/systemd/piglow-log-bridge.user.service.example
# PiGlow log bridge — drives PiGlow from OpenClaw app log (run start/done)
# Copy to: ~/.config/systemd/user/piglow-log-bridge.service
# Then: systemctl --user daemon-reload && systemctl --user enable piglow-log-bridge.service && systemctl --user start piglow-log-bridge.service
# Logs: journalctl --user -u piglow-log-bridge.service -f
# Requires: PiGlow service running (piglow-service.service), jq, agent logs at agents/AGENT_NAME/logs/
# Edit WorkingDirectory and the script argument to your repo path and agent name.

[Unit]
Description=PiGlow log bridge for pibot (activity on run start/done)
After=network-online.target piglow-service.service
Wants=piglow-service.service

[Service]
Type=simple
WorkingDirectory=/path/to/ironclaw
ExecStart=/path/to/ironclaw/scripts/piglow-log-bridge.sh AGENT_NAME
Restart=on-failure
RestartSec=10

[Install]
WantedBy=default.target
