# scripts/systemd/rfid-notify-bridge.user.service.example
# RFID notify bridge — runs the poll script once (invoked by the timer)
# Copy to: ~/.config/systemd/user/rfid-notify-bridge.service
# Then: systemctl --user daemon-reload && systemctl --user enable rfid-notify-bridge.timer && systemctl --user start rfid-notify-bridge.timer
# Timer runs every 10s so the agent is notified even if the daemon's immediate gateway call failed.

[Unit]
Description=RFID notify bridge (one-shot, called by timer)

[Service]
Type=oneshot
WorkingDirectory=/path/to/ironclaw
EnvironmentFile=/path/to/ironclaw/agents/AGENT_NAME/.env
ExecStart=/usr/bin/python3 /path/to/ironclaw/agents/AGENT_NAME/workspace/rfid/rfid_notify_bridge.py
